# recipes-support/syslog-ng/files/powerstate_update.path
[Unit]
Description=Power State Update

[Path]
PathChanged=/tmp/.standby
Unit=powerstate_update.service

[Install]
WantedBy=multi-user.target

